# .config/systemd/user/cuthulu.service
[Unit]
Description=Cuthulu client
After=network-online.target graphical-session.target
Wants=network-online.target graphical-session.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
ExecStart=%h/.local/bin/cuthulu
Environment=HOME=%h
Environment=DISPLAY=:0
Environment=CUTHULU_SSH_HOST=thinkpad.holdenitdown.net
Environment=CUTHULU_SSH_USER=rfhold
Environment=CUTHULU_NO_STUBS=1
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
