# .config/systemd/user/tmux-agent.service
[Unit]
Description=Tmux session agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/go/bin/go run %h/dot/cmd/tmux-agent
Environment=HOME=%h
Environment=GOPATH=%h/go
Environment=PATH=/usr/local/go/bin:%h/go/bin:/usr/local/bin:/usr/bin:/bin
Environment=SERVER_URL=https://tmux.holdenitdown.net
Environment=SSH_HOST=thinkpad.holdenitdown.net
Environment=SSH_USER=rfhold
Environment=SKIP_AUTH=true
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
