# etc/systemd/user/cuthulu.service
[Unit]
Description=Cuthulu OpenCodes client
After=network-online.target graphical-session.target
Wants=network-online.target graphical-session.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
ExecStart=%h/.local/bin/cuthulu
Environment=HOME=%h
Environment=DISPLAY=:0
Environment=OPENCODES_SSH_HOST={{ hostname }}.holdenitdown.net
Environment=OPENCODES_SSH_USER=rfhold
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
