# etc/launchagents/cuthulu.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0"><dict>
  <key>Label</key><string>dev.rholden.cuthulu</string>
  <key>ProgramArguments</key>
  <array><string>{{ app_path }}</string></array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key><string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>CUTHULU_SSH_HOST</key><string>macbook.holdenitdown.net</string>
    <key>CUTHULU_SSH_USER</key><string>rfhold</string>
    <key>SSH_AUTH_SOCK</key><string>/Users/rfhold/.gnupg/S.gpg-agent.ssh</string>
  </dict>
  <key>StandardOutPath</key><string>/tmp/cuthulu.log</string>
  <key>StandardErrorPath</key><string>/tmp/cuthulu.log</string>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
</dict></plist>
