# etc/launchagents/opencodes-tray.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0"><dict>
  <key>Label</key><string>dev.rfholden.opencodes-tray</string>
  <key>ProgramArguments</key>
  <array><string>{{ tray_path }}</string></array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key><string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>OPENCODES_SSH_HOST</key><string>macbook.holdenitdown.net</string>
    <key>OPENCODES_SSH_USER</key><string>rfhold</string>
  </dict>
  <key>StandardOutPath</key><string>/tmp/opencodes-tray.log</string>
  <key>StandardErrorPath</key><string>/tmp/opencodes-tray.log</string>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
</dict></plist>
